# examples/crontab
# m    h  dom mon dow   command
*/30   *   *   *   *    $HOME/report-aggregator/cron_job.sh
